# kindle_hid_passthrough/hid-passthrough-dev.upstart
# HID Passthrough - Bluetooth HID host with UHID passthrough
# For development with Python environment
start on started wmt
stop on stopping filesystems

respawn
respawn limit 5 60

env KINDLE_HID_BASE=/mnt/us/kindle_hid_passthrough

script
    # Load Bluetooth kernel module if not loaded
    if ! lsmod | grep -q wmt_cdev_bt; then
        /sbin/insmod /lib/modules/4.9.77-lab126/extra/wmt_cdev_bt.ko || true
        sleep 1
    fi

    # Kill any conflicting processes
    killall bluetoothd 2>/dev/null || true
    killall vhci_stpbt_bridge 2>/dev/null || true
    sleep 1

    exec /mnt/us/python3.10-kindle/python3-wrapper.sh /mnt/us/kindle_hid_passthrough/daemon.py
end script
