# kindle_hid_passthrough/hid-passthrough-dev.upstart
# HID Passthrough - Bluetooth HID host with UHID passthrough
# Requires WMT (Bluetooth hardware) to be ready
start on started wmt
stop on stopping filesystems

respawn
respawn limit 5 60

env KINDLE_HID_BASE=/mnt/us/kindle_hid_passthrough

exec /mnt/us/python3.10-kindle/python3-wrapper.sh /mnt/us/kindle_hid_passthrough/daemon.py
